# backend/src/wb-clusters/com.wb-automation.query-map-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent — еженедельное обновление карты запросов кластеров.

  Что делает каждое воскресенье в 10:00 по местному времени:
    1. Открывает SSH-туннель к production PostgreSQL (порт 15432).
    2. Загружает все уникальные advertId из wb_campaign_products.
    3. Открывает единственное окно Safari на cmp.wildberries.ru.
    4. Батчами по 40 advertId запускает параллельные fetch() в браузере
       (concurrency=10), скачивая /api/v5/words-clusters для каждой РК.
    5. Сохраняет результаты в wb_cabinet_cluster_queries через UNNEST.
    6. Закрывает SSH-туннель.

  Предварительные требования:
    - Safari должен быть открыт с авторизованной сессией cmp.wildberries.ru.
      Если начнут появляться ошибки auth-error — войти вручную.
    - SSH-ключ должен разрешать вход без пароля на root@95.163.226.154.

  Установка:
    cp com.wb-automation.query-map-update.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.wb-automation.query-map-update.plist

  Запустить немедленно (без ожидания воскресенья):
    launchctl start com.wb-automation.query-map-update

  Отключить:
    launchctl unload ~/Library/LaunchAgents/com.wb-automation.query-map-update.plist

  Логи:
    tail -f ~/Downloads/wb-query-map-update.log
    tail -f ~/Downloads/wb-query-map-update-error.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wb-automation.query-map-update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/elizavetabojkova/Downloads/Маркетплейс/backend/src/wb-clusters/run-query-map-update.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/elizavetabojkova/Downloads/Маркетплейс/backend</string>

    <!-- Каждое воскресенье в 10:00 по местному времени (Weekday=0 = воскресенье) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/elizavetabojkova/Downloads/wb-query-map-update.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/elizavetabojkova/Downloads/wb-query-map-update-error.log</string>

    <!-- Не запускать при загрузке, только по расписанию -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
