# backend/src/wb-clusters/com.wb-automation.monthly-frequency-import.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent for daily WB search analytics frequency import.

  What it does every day at 09:00 local time:
    1. Opens SSH tunnel to production PostgreSQL (port 15432).
    2. Finds or creates a SUCCESS report on WB servers via seller-content API.
    3. Opens the Downloads Manager panel in Safari and clicks "Скачать".
       The portal's React handler fetches the file with httpOnly session cookies
       and saves it to ~/Downloads as a ZIP.
    4. Extracts the XLSX from the ZIP, parses ~300k query-frequency rows.
    5. Replaces the wb_search_query_frequencies snapshot in PostgreSQL.
    6. Closes the SSH tunnel.

  Prerequisites:
    - Safari must be open with seller.wildberries.ru logged in (fresh session).
      Re-login manually if downloads start failing with 403.
    - SSH key must allow passwordless login to root@95.163.226.154.

  To install:
    cp com.wb-automation.monthly-frequency-import.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.wb-automation.monthly-frequency-import.plist

  To test immediately (without waiting for 09:00):
    launchctl start com.wb-automation.monthly-frequency-import

  To unload:
    launchctl unload ~/Library/LaunchAgents/com.wb-automation.monthly-frequency-import.plist

  Logs:
    tail -f ~/Downloads/wb-frequency-import.log
    tail -f ~/Downloads/wb-frequency-import-error.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wb-automation.monthly-frequency-import</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/elizavetabojkova/wb-automation/backend/src/wb-clusters/run-monthly-frequency-import.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/elizavetabojkova/wb-automation/backend</string>

    <!-- Run every day at 09:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/elizavetabojkova/wb-automation/logs/wb-frequency-import.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/elizavetabojkova/wb-automation/logs/wb-frequency-import-error.log</string>

    <!-- Do not run on load, only at the scheduled time -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
